# debian/postinst
depmod -a
chmod 0644 /lib/systemd/system/marvell_teralynx_driver.service
systemctl enable marvell_teralynx_driver.service
systemctl start marvell_teralynx_driver.service
